# crontab
CRON_TZ=America/Chicago
0 7 * * 1 python /app/sync.py
