# crontab
CRON_TZ=America/Chicago
0 1 * * 3 python /app/fab_local_dfw_events.py
0 2 * * 3 python /app/fab_major_global_events.py
